# etc/xdg/autostart/env-display.desktop
[Desktop Entry]
Type=Application
Name=env_display
Comment=Set display variable for systemd user services
NoDisplay=false
Exec=/usr/bin/dbus-update-activation-environment --systemd DISPLAY WAYLAND_DISPLAY XDG_CURRENT_DESKTOP=Wayfire
