# etc/xdg/autostart/xsettingsd.desktop
[Desktop Entry]
Name=Xsettings daemon for Xwayland
Exec=wfrespawn xsettingsd
Terminal=false
Type=Application
NoDisplay=true
OnlyShowIn=wayfire
